# docker/root/etc/s6/gitbundle/run
#!/bin/bash
[[ -f ./setup ]] && source ./setup

pushd /app/gitbundle >/dev/null
exec su-exec $USER /usr/local/bin/gitbundle web
popd
